# modules/share-extra/airootfs.any/etc/profile.d/disable-beep.sh
#!/usr/bin/env bash

# disable-beep.sh
#
# Script to disable PC beep sound
#

setterm -blength 0
